# root/etc/services.d/gazee/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python3 /app/gazee/Gazee.py \
	-c /config
